# router/conf/dyndns/dyndns.timer
[Unit]
Description=Automatically update dynamic dns

[Timer]
OnBootSec=2min
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target