# router/conf/dyndns/dyndns.service
[Unit]
Description="Digital Ocean dynamic DNS"

[Service]
ExecStart=/usr/local/sbin/dyndns/dyndns_updater.py
EOF