# router/conf/borg/borg.service
[Unit]
Description="Automate backups w/ borg"

[Service]
EnvironmentFile=/usr/local/sbin/configs/.env
ExecStart=/usr/local/sbin/backup.sh
